# Actividad4/.service
[Unit]
Description=Saludo diario
After=network-online.target

[Service]
ExecStart=/usr/local/bin/saludo.sh

[Install]
WantedBy=multi-user.target
